# taos-miner.service
[Unit]
Description=TAOS Bittensor Miner Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/sn79-miner
ExecStart=/home/dev/sn79-miner/run_miner_service.sh
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taos-miner

# Security settings
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/dev/sn79-miner/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
